# scripts/drummachine.service
[Unit]
Description=Raspberry Pi Drum Machine v2.0
After=multi-user.target sound.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/DRUMMACHINE
ExecStart=/usr/bin/python3 /home/pi/DRUMMACHINE/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Variables de entorno necesarias
Environment="DISPLAY=:0"
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

